# docker/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=false
logfile=/var/www/storage/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ---------------------------------------------------------------------------
# Nginx - Servidor web
# ---------------------------------------------------------------------------
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
priority=10
stdout_logfile=/var/www/storage/logs/nginx-out.log
stderr_logfile=/var/www/storage/logs/nginx-err.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

# ---------------------------------------------------------------------------
# PHP-FPM - FastCGI Process Manager
# ---------------------------------------------------------------------------
[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
startretries=5
priority=5
stdout_logfile=/var/www/storage/logs/php-fpm-out.log
stderr_logfile=/var/www/storage/logs/php-fpm-err.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

# ---------------------------------------------------------------------------
# Laravel Reverb - Servidor WebSocket
# ---------------------------------------------------------------------------
[program:reverb]
command=php /var/www/artisan reverb:start --host=0.0.0.0 --port=8080 --no-interaction
directory=/var/www
user=www-data
autostart=true
autorestart=true
startretries=10
startsecs=3
priority=20
stdout_logfile=/var/www/storage/logs/reverb-out.log
stderr_logfile=/var/www/storage/logs/reverb-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=5MB

# ---------------------------------------------------------------------------
# Laravel Queue Worker - Procesador de trabajos en cola
# ---------------------------------------------------------------------------
[program:queue-worker]
command=php /var/www/artisan queue:work --sleep=3 --tries=3 --timeout=90 --max-time=3600 --no-interaction
directory=/var/www
user=www-data
autostart=true
autorestart=true
startretries=10
startsecs=2
priority=30
numprocs=2
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/www/storage/logs/queue-out.log
stderr_logfile=/var/www/storage/logs/queue-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=5MB

# ---------------------------------------------------------------------------
# Laravel Scheduler - Tareas programadas (reemplaza cron)
# ---------------------------------------------------------------------------
[program:scheduler]
command=php /var/www/artisan schedule:work --no-interaction
directory=/var/www
user=www-data
autostart=true
autorestart=true
startretries=5
priority=40
stdout_logfile=/var/www/storage/logs/scheduler-out.log
stderr_logfile=/var/www/storage/logs/scheduler-err.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
